# crontab
# START CRON JOB
0 * * * * /usr/local/bin/python3 /app/cfddns.py YourDomainHere APIKeyHere > /proc/1/fd/1 2>/proc/1/fd/2
# END CRON JOB
